# raspiframe-kiosk.service
# このファイルは setup_dlna.sh によって自動生成されます
# 直接編集しないでください
# This file is automatically generated by setup_dlna.sh
# Do not edit directly

[Unit]
Description=Raspiframe Kiosk Mode Startup
After=raspiframe.service graphical.target
Wants=raspiframe.service

[Service]
Type=oneshot
User=pi
Environment="DISPLAY=:0"
Environment="XAUTHORITY=/home/pi/.Xauthority"
ExecStart=/home/jd/raspiframe/startup_pipeline.sh
RemainAfterExit=yes

[Install]
WantedBy=graphical.target



